# overlays/CARAMBOLA2/c/etc/crontab
#minute	hour	mday	month	wday	who	command
# until we can run a lightweight ntp client in the backgronud..
0/15	*	*	*	*	root	/etc/rc.d/ntpdate onestart
